# multi_agent_meeting/deploy/multi-agent-meeting.service
[Unit]
Description=Multi Agent Meeting System
Documentation=https://github.com/yinjianhui/camel_ai
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/camel_ai
Environment=PATH=/opt/camel_ai/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/opt/camel_ai
EnvironmentFile=/opt/camel_ai/multi_agent_meeting/backend/.env
ExecStart=/opt/camel_ai/venv/bin/python multi_agent_meeting/backend/app_new.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=multi-agent-meeting

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/camel_ai
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
